# app/conf.d/xsession.conf
[program:xsession]
command=/etc/X11/Xsession
user=developer
autorestart=true
stdout_logfile=/app/xsession.log
redirect_stderr=true