# image/app/conf.d/start.conf
[program:start]
command=/app/start.sh
autorestart=true
stdout_logfile=/app/log/start.log
redirect_stderr=true